# config/storaged.service
[Unit]
Description=CedarGraph Storage Server (storaged)
Documentation=https://docs.cedargraph.io
After=network.target
Wants=network.target

[Service]
Type=simple
User=cedar
Group=cedar

# Working directory
WorkingDirectory=/var/lib/cedar/storage

# Executable
ExecStart=/usr/local/bin/storaged --config /etc/cedar/storaged.conf
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy - auto recovery on failure
Restart=on-failure
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment
Environment="RUST_BACKTRACE=1"
Environment="CEDAR_LOG_LEVEL=info"

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/cedar/storage /var/log/cedar

# Logging
StandardOutput=append:/var/log/cedar/storaged.log
StandardError=append:/var/log/cedar/storaged.error.log
SyslogIdentifier=storaged

[Install]
WantedBy=multi-user.target
